# .config/launchd/com.miniminjae.codex-account-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.miniminjae.codex-account-usage</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/miniminjae/.local/bin/codex-account-usage</string>
    <string>monitor</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>/Users/miniminjae/.local/state/codex-account-usage/monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/miniminjae/.local/state/codex-account-usage/monitor.error.log</string>
</dict>
</plist>
