# .config/launchd/com.miniminjae.agent-notify-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.miniminjae.agent-notify-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/miniminjae/.local/bin/agent-notify</string>
    <string>sweep</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/miniminjae/.local/state/agent-notify/sweep.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/miniminjae/.local/state/agent-notify/sweep.error.log</string>
</dict>
</plist>
